# src/.config/systemd/user/swww-daemon.service
[Unit]
Description=SWWW Daemon

[Service]
ExecStart=killall swww-daemon; /usr/bin/swww-daemon
Restart=on-failure

[Install]
WantedBy=default.target
